# online-desktop/trunk/autostart/bigboard.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Sidebar
Exec=bigboard
Version=1.0
X-GNOME-Autostart-enabled=true
