# online-desktop/trunk/autostart/browser.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Browser
# FIXME - redhat/fedora specific?
Exec=htmlview
Version=1.0
X-GNOME-Autostart-enabled=true
