# online-desktop/trunk/autostart/weblogindriver.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Web Login Driver
Exec=web-login-driver
Type=Application
Version=1.0
X-GNOME-Autostart-enabled=true

